# deploy/bot.service
[Unit]
Description=Telegram Assignment Bot
After=network.target

[Service]
User=botuser
Group=botuser
WorkingDirectory=/opt/tbtbfx-create
Environment=BOT_TOKEN=
Environment=CHANNEL_ID=
Environment=ADMIN_ID=
ExecStart=/usr/bin/python3 /opt/tbtbfx-create/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
